# docker_smbldap/docker_samba/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true

[program:smbd]
command=/usr/sbin/smbd -F --no-process-group
autorestart=true
